# view_updater/view_updater.cron
0 * * * * root /var/lounge/bin/view_updater.py
